# common/updown-badge-server.service
[Unit]
Description=Service for generating badges for updown.io checks
Documentation=https://kaos.sh/updown-badge-server
Requires=network.target remote-fs.target
After=network.target remote-fs.target

[Service]
ExecStart=/usr/bin/updown-badge-server --no-color
ExecReload=kill -HUP $MAINPID
SyslogIdentifier=updown-badge-server

[Install]
WantedBy=multi-user.target
